# wiresheet/rootfs/etc/s6-overlay/s6-rc.d/wiresheet-server/run
#!/command/with-contenv bashio

echo "=== WIRESHEET SERVER START ==="

cd /app
exec node server/index.js
