# wiresheet/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio

echo "Starte Nginx..."
exec nginx -g "daemon off;"
